# systemd/user/picom.service
[Unit]
Description=Picom daemon

[Service]
Type=simple
ExecStart=/usr/bin/picom


[Install]
WantedBy=graphical.target
